# supervisor-app.conf
[supervisord]
nodaemon=true

[program:mqttd]
command = /home/file/emqttd/bin/emqttd start

